# ems-kalibrierung.service
[Unit]
Description=EMS Vorhersage-Kalibrierung (PV Solcast + Hausverbrauch)
After=network-online.target influxdb.service
Wants=network-online.target

[Service]
Type=oneshot
User=ems
Group=ems
WorkingDirectory=/opt/ems
# Erzeugt kalibrierung.yaml (Report) und kalibrierung_profil.yaml (PV-Profil).
ExecStart=/opt/ems/.venv/bin/python -m kalibrierung --config /opt/ems/config.yaml --lookback-days 730 --test-days 365
Environment=PYTHONUNBUFFERED=1
# Haertung analog ems.service (kein Port noetig -> keine Capabilities).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ems
PrivateTmp=true
