# ems-history.service
[Unit]
Description=EMS E3DC-Historie-Import (RSCP -> lokale SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ems
Group=ems
WorkingDirectory=/opt/ems
# Holt die letzten Tage in 15-min-Auflösung nach (idempotent, UPSERT) - fängt
# auch Lücken, wenn der Pi aus war. 2 Tage = ~192 RSCP-Aufrufe. Grosse
# Erstbefüllung einmalig manuell: rscp_import.py --resolution 15min --days N.
ExecStart=/opt/ems/.venv/bin/python rscp_import.py --config /opt/ems/config.yaml --resolution 15min --days 2
Environment=PYTHONUNBUFFERED=1
# Haertung analog ems.service (nur /opt/ems beschreibbar).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ems
PrivateTmp=true
